# .docker/cronjobs/crontabs
* * * * * root php /usr/share/nginx/artisan schedule:run >/var/log/gustavoppk 2>&1