# .docker/cronjobs/crontab
* * * * * root php /usr/share/nginx/artisan schedule:run >> /var/log/cron.log 2>&1